# systemd/radxa-network-monitor.service
[Unit]
Description=Radxa Network Monitor Service
Documentation=https://github.com/SMT03/new-wave-linux
After=network-online.target radxa-manager.service
Wants=network-online.target
BindsTo=radxa-manager.service

[Service]
Type=simple
User=radxa
Group=radxa
WorkingDirectory=/opt/radxa
Environment=PYTHONPATH=/opt/radxa
ExecStart=/opt/radxa/venv/bin/python /opt/radxa/network_monitor.py --continuous --config /etc/radxa/settings.yaml
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/radxa /var/lib/radxa /tmp/radxa
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true

# Resource limits
MemoryMax=256M
CPUQuota=50%
TasksMax=20

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=radxa-network-monitor

[Install]
WantedBy=multi-user.target
